# github_hooks_dev/deployment_installer/run_pyelixys_server/dependencies/config/autostart/startup_script.desktop
[Desktop Entry]
Type=Application
Name=Elixys Server
Exec=gnome-terminal -e /usr/sbin/elixys
Comment=Starts Elixys
Terminal=False
Type=Application
Categories=Utility;Application;